# .config/systemd/user/hyprshade.service
[Unit]
Description=Apply screen filter

[Service]
Type=oneshot
ExecStart=/usr/bin/hyprshade auto
