# .config/systemd/user/hyprshade.timer
[Unit]
Description=Apply screen filter on schedule

[Timer]
OnCalendar=*-*-* 06:00:00
OnCalendar=*-*-* 20:00:00

[Install]
WantedBy=timers.target
