# .config/systemd/user/hypridle.service
[Unit]
Description=Start hypridle daemon
After=wayland.target

[Service]
ExecStart=hypridle

[Install]
WantedBy=wayland.target
