# .config/systemd/user/user-graphical-session.target
[Unit]
Description=Shared wayland-xorg graphical session target
Requires=default.target
