# .config/systemd/user/wayland.target
[Unit]
Description=Wayland session target
Requires=user-graphical-session.target
Conflicts=xorg.target
