# .config/systemd/user/xorg.target
[Unit]
Description=x11-xorg user session target
Requires=user-graphical-session.target
Conflicts=wayland.target
