# disk-root/etc/systemd/system/getty@tty1.service.d/getty@tty1.service.d/vlk-login.conf
[Service]
ExecStartPre=/bin/sh -c 'setleds -D +num < /dev/%I'
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- vlk' --noclear --skip-login - $TERM
